# src/rpi-zero-usb-iso.service
[Unit]

[Service]
ExecStart=%h/.local/bin/rpi-zero-usb-iso iso
Type=oneshot

[Install]
WantedBy=default.target
